# Dockerfile/elk/resources/elasticsearch.conf
[program:elasticsearch]
command=sudo -u elasticsearch JAVA_HOME=/opt/jdk/jre /opt/elasticsearch/bin/elasticsearch -p /var/run/elasticsearch/elasticsearch.pid -Des.default.path.home=/opt/elasticsearch -Des.default.path.logs=/opt/elasticsearch/logs -Des.default.path.work=/tmp/elasticsearch -Des.default.path.conf=/opt/elasticsearch/config -Des.default.path.repo=/opt/elasticsearch/snapshot
