# Dockerfile/elk/resources/logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash -f /opt/logstash/data_report.conf
stdout_logfile=/var/log/logstash.log
redirect_stderr=true
